# vm_setup/99-in_changed.rules
ACTION=="change", RUN+="/home/gladiator/in_change_handler.sh"
