# vm_setup/in_changed@.service
[Unit]
Description=Change in %i

[Service]
Type=simple
ExecStart=/home/gladiator/in_changed.sh %i
